# novel_frame/deploy/novel-agent.service
[Unit]
Description=AI Novel Agent Backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data

# >>> 修改为你的实际项目路径 <<<
WorkingDirectory=/your/path/to/backend
ExecStart=/your/path/to/backend/.venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000 --workers 1

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

Environment=PYTHONUNBUFFERED=1

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# >>> 修改为你的实际 data 路径 <<<
ReadWritePaths=/your/path/to/backend/data

[Install]
WantedBy=multi-user.target
